# web/deploy/supervisord-user.conf
[unix_http_server]
file=/var/run/supervisord-user.sock

[supervisorctl]
serverurl=unix:///var/run/supervisord-user.sock

[supervisord]
logfile=/dev/shm/supervisord_user.log
nodaemon=false
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:resonance-api]
command=bash -c 'cd /opt/resonance-api && RESONANCE_BINARY=/opt/resonance/build/acoustic_sim python3 -m uvicorn server:app --host 127.0.0.1 --port 8420 --workers 2'
directory=/opt/resonance-api
autostart=true
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/var/log/resonance-api.log
stderr_logfile=/var/log/resonance-api.log
stopasgroup=true
killasgroup=true
